# building/build-debs/homeworld-keysystem/services/update-keyclient-config.timer
[Unit]
Description=homeworld keyclient configuration update

[Timer]
OnBootSec=15sec
OnUnitActiveSec=30sec

[Install]
WantedBy=timers.target
