# building/build-debs/homeworld-keysystem/services/update-keyclient-config.service
[Unit]
Description=homeworld keyclient configuration update

[Service]
ExecStart=/usr/lib/homeworld/update-keyclient-config.sh
Restart=never
