# building/src/setup-queue/setup-queue.timer
[Unit]
Description=homeworld autodeploy queue

[Timer]
OnBootSec=15sec
OnUnitActiveSec=15sec

[Install]
WantedBy=timers.target
