# platform/pull-monitor/pull-monitor/pull-monitor.service.in
[Unit]
Description=Container pull monitor
Requires=network-online.target
After=network-online.target

[Service]
Slice=machine.slice
ExecStart=/usr/bin/pull-monitor docker://homeworld.private/pullcheck@{PULLCHECK_DIGEST}
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
