# deployment/services/flannel.service
[Unit]
Description=Network fabric for containers
After=etcd.service

[Service]
Type=notify
Restart=always
Environment="NOTIFY_SOCKET=/run/systemd/notify"

ExecStart=/usr/lib/hyades/launch-flannel.sh

[Install]
WantedBy=multi-user.target
