# building/build-debs/homeworld-rkt/services/rkt-gc.timer
[Unit]
Description=Periodic Garbage Collection for rkt

[Timer]
OnActiveSec=0s
OnUnitActiveSec=5m

[Install]
WantedBy=multi-user.target
