# building/build-debs/homeworld-rkt/services/rkt-gc.service
[Unit]
Description=Garbage Collection for rkt

[Service]
Type=oneshot
ExecStart=/usr/bin/rkt gc --grace-period=15m
