# platform/aci-pull-monitor/aci-pull-monitor/aci-pull-monitor.service.in
[Unit]
Description=aci pull monitor
Requires=network-online.target
After=network-online.target

[Service]
Slice=machine.slice
ExecStart=/usr/bin/aci-pull-monitor docker://homeworld.private/pullcheck:{PULLCHECK_VERSION}
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
