# platform/setup-queue/setup-queue.service
[Unit]
Description=homeworld autodeploy queue
After=keyclient.service

[Service]
Type=oneshot
ExecStart=/usr/bin/setup-queue-deploy
Restart=no
